# Startup_Files/init_can.service
[Unit]
Description=CAN Init Script

[Service]
ExecStart=/home/nvidia/Palletfinder/Startup_Files/init_can.sh

[Install]
WantedBy=multi-user.target