# Startup_Files/init_docker.service
[Unit]
Description=Docker Init ...

[Service]
ExecStart=/home/nvidia/Palletfinder/Startup_Files/Docker_start.sh

[Install]
WantedBy=multi-user.target